# utils/climate/ha-climate-serve@ha.service
[Unit]
Description=store mqtt data in db
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/ha/home-assistant-setup/utils/climate
ExecStart=/usr/bin/python3 /home/ha/home-assistant-setup/utils/climate/serve.py

[Install]
WantedBy=multi-user.target
